# docker/cron/crontab
35 * * * * /opt/console -option=squad > /proc/1/fd/1 2> /proc/1/fd/2
