# docker/cron/crontab
30 * * * * /opt/console -option=team > /proc/1/fd/1 2> /proc/1/fd/2
